# qca/src/recipes-qcawifi-oss/initoverlay/files/initoverlay/etc/profile.d/overlay.sh
# Change the directory once the overlay has been mounted
cd /
